# read_iops.service
[Service]
Slice=io_parent.slice
IOAccounting=yes
ExecStart=/root/read_iops.sh
